# etc/cron.d/rtl-wx
2 */12 * * * root /opt/rtl-wx/bin/accu >/var/log/accu.log 2>&1
