# opt/rtl-wx/DEBIAN/postinst
#!/bin/bash
a2ensite wx
systemctl restart apache2
systemctl daemon-reload
systemctl enable wx
systemctl restart wx
